# deploy/steam-bot.service
[Unit]
Description=Steam Inventory Monitor Telegram Bot
After=network.target

[Service]
Type=simple
User=steam-bot
EnvironmentFile=/opt/steam-inventory-bot/.env
WorkingDirectory=/opt/steam-inventory-bot
ExecStart=/usr/bin/python3 /opt/steam-inventory-bot/main.py
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=steam-bot

[Install]
WantedBy=multi-user.target
